# deb/common/docker-ce.preinst
#!/bin/sh
set -e

case "$1" in
	upgrade)
		invoke-rc.d docker stop || exit $?
		;;
	abort-*)
		# How'd we get here??
		exit 1
		;;
	*)
		;;
esac

#DEBHELPER#
